# contrib/burrow.service
# systemd unit for the burrow daemon (system-wide install).
#
# The .deb/.rpm packages install this unit but leave it DISABLED, because
# burrow cannot start until it has been initialised once by a human:
#
#   sudo burrowctl init      # generates the repo key, PRINTS your 24-word
#                            # recovery phrase — store it OFF this machine
#   sudo systemctl enable --now burrow
#
# State layout (via the BURROW_* env vars below):
#   /etc/burrow      config.toml + repo.key  (secret; dir 0750, key 0600)
#   /var/lib/burrow  metadata DB + blob store
#   /run/burrow      CLI control socket
#
# NixOS users: prefer the flake's nixosModules.burrow instead of this unit.

[Unit]
Description=burrow — distributed backup among friends
Documentation=https://github.com/solidsilver/burrow
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=burrow
Group=burrow
Environment=BURROW_CONFIG_DIR=/etc/burrow
Environment=BURROW_DATA_DIR=/var/lib/burrow
Environment=BURROW_SOCKET=/run/burrow/daemon.sock
ExecStart=/usr/bin/burrow daemon run
Restart=on-failure
RestartSec=5
TimeoutStopSec=15

# systemd-managed directories: created with the right owner on start, so no
# maintainer script has to mkdir/chown them.
StateDirectory=burrow
StateDirectoryMode=0700
ConfigDirectory=burrow
ConfigDirectoryMode=0750
RuntimeDirectory=burrow
RuntimeDirectoryMode=0700

# --- Sandboxing (see `systemd-analyze security burrow`) ---
NoNewPrivileges=true
ProtectSystem=strict
# read-only, NOT true: backups must be able to read /home (every example
# backs up /home/...). Writes under /home still fail — ProtectSystem=strict
# already confines writes to the State/Config/Runtime directories; a restore
# to other paths needs a ReadWritePaths drop-in, same as before.
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
# AF_NETLINK is required: iroh enumerates local interfaces to gather QUIC
# connection candidates. AF_INET/AF_INET6 for QUIC, AF_UNIX for the CLI socket.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077

[Install]
WantedBy=multi-user.target
